# scripts/linux/systemd/appsuite-itsm-dev.service
[Unit]
Description=AppSuite ITSM Management System - Development Environment
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=kensan
Group=kensan
WorkingDirectory=/mnt/LinuxHDD/Appsuite-ITSM-Management/WebUI-Production
Environment="NODE_ENV=development"
Environment="PORT=3100"

# http-serverを起動（開発環境: ポート3100固定）
ExecStart=/usr/bin/npx http-server -p 3100 -c-1 --cors

# 自動再起動設定
Restart=on-failure
RestartSec=10

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=appsuite-itsm-dev

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
